# esl-random/s3_watchdog_service.service
[Unit]
Description=Python S3 Watchdog Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /path/to/s3_watchdog_service.py <bucket_name> <dir1> [<dir2> ...]
Restart=on-failure
User=ec2-user

[Install]
WantedBy=multi-user.target
